# BashScripts/google-chrome-scripts/google-config-backup.timer
[Unit]
Description=Google Chrome backup timer

[Timer]
OnUnitActiveSec=10m
OnBootSec=0s

[Install]
WantedBy=timers.target
