# BashScripts/firefox-scripts/google-config-backup.timer
[Unit]
Description=Firefox backup timer

[Timer]
OnUnitActiveSec=30m
OnBootSec=30m

[Install]
WantedBy=timers.target
